# config/paste.supervisor.ini
[program:paste]
directory=/home/app/deployed/paste
command=gunicorn -b unix:/tmp/gunicorn.paste.sock app.app:app
user=app
environment=DATABASE=sqlite:////home/app/deployed/paste.db
